# src/Lyfe/Lyfe.service
[Unit]
Description=Lyfe

[Service]
User=pi
Environment=ASPNETCORE_URLS=http://*:5000
WorkingDirectory=/home/pi/Lyfe
ExecStart=/home/pi/.dotnet/dotnet Lyfe.dll
Restart=on-failure

[Install]
WantedBy=multi-user.target
